# com.williamcalhoun.spotify-taste-graph.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamcalhoun.spotify-taste-graph</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/calhoun_r2/Desktop/development/spotify-taste-graph</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/calhoun_r2/Library/Logs/spotify-taste-graph.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/calhoun_r2/Library/Logs/spotify-taste-graph.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
